# hadoop-tools/hadoop-docker/src/main/wrappers/etc/profile.d/hadoop.sh
export HADOOP_HOME=/opt/hadoop-${HADOOP_VERSION}
export HADOOP_CONF_DIR=/etc/hadoop
export JAVA_HOME=$(find /usr/lib/jvm/* -prune | tail -1)
